# rc.local
#!/bin/bash
rfkill block bluetooth
exit 0

